# WRX36/bin/extras/files/etc/uci-defaults/zz_custom_unbound
mkdir -p /etc/unbound/root-anchors
sudo -u unbound unbound-control-setup -d /etc/unbound
wget https://data.iana.org/root-anchors/icannbundle.pem -O /etc/unbound/root-anchors/icannbundle.pem
wget https://data.iana.org/root-anchors/root-anchors.p7s -O /etc/unbound/root-anchors/root-anchors.p7s
wget https://data.iana.org/root-anchors/root-anchors.xml -O /etc/unbound/root-anchors/root-anchors.xml
wget https://www.internic.net/domain/named.root -O /etc/unbound/root.hints
cp /var/lib/unbound/root.key /etc/unbound
sudo -u unbound unbound-anchor -a /etc/unbound/root.key -c /etc/unbound/root-anchors/icannbundle.pem -x /etc/unbound/root-anchors/root-anchors.xml -s  /etc/unbound/root-anchors/root-anchors.p7s -r /etc/unbound/root.hints
chown -R unbound:unbound /etc/unbound
service unbound restart
