# WRX36/bin/extra/files/etc/uci-defaults/z8-chrony
#!/bin/sh

{ [ -f /etc/config/chrony ] && sed s/'[[:space:]]'//g /etc/config/chrony | grep -qE '.+' && [ -f /rom/etc/config/chrony ] && sed s/'[[:space:]]'//g /rom/etc/config/chrony | grep -qE '.+' && ! [ "$(cat /etc/config/chrony)" = "$(cat /rom/etc/config/chrony)" ]; } || {

touch /etc/config/chrony

uci batch <<'EOF'
delete chrony.@pool[0]
add chrony pool
set chrony.@pool[0]=pool
set chrony.@pool[0].hostname='stratum1.time.cifelli.xyz'
set chrony.@pool[0].maxpoll='12'
set chrony.@pool[0].iburst='1'
set chrony.@pool[0].nts='1'

delete chrony.@dhcp_ntp_server[0]
add chrony dhcp_ntp_server
set chrony.@dhcp_ntp_server[0]=dhcp_ntp_server
set chrony.@dhcp_ntp_server[0].iburst='1'
set chrony.@dhcp_ntp_server[0].disabled='0'

delete chrony.@allow[0]
add chrony allow
set chrony.@allow[0]=allow
set chrony.@allow[0].interface='lan'

delete chrony.@makestep[0]
add chrony makestep
set chrony.@makestep[0]=makestep
set chrony.@makestep[0].threshold='1.0'
set chrony.@makestep[0].limit='3'

delete chrony.@nts[0]
add chrony nts
set chrony.@nts[0]=nts
set chrony.@nts[0].rtccheck='1'
set chrony.@nts[0].systemcerts='1'
EOF
uci commit chrony
}
exit 0
