# WRX36/bin/extra/files/etc/uci-defaults/z8-fstab
#!/bin/sh

{ [ -f /etc/config/fstab ] && sed s/'[[:space:]]'//g /etc/config/fstab | grep -qE '.+' && [ -f /rom/etc/config/fstab ] && sed s/'[[:space:]]'//g /rom/etc/config/fstab | grep -qE '.+' && ! [ "$(cat /etc/config/fstab)" = "$(cat /rom/etc/config/fstab)" ]; } || {

touch /etc/config/fstab

uci batch <<'EOF'
delete fstab.@global[0]
add fstab global
set fstab.@global[0]=global
set fstab.@global[0].anon_swap='0'
set fstab.@global[0].anon_mount='0'
set fstab.@global[0].auto_swap='1'
set fstab.@global[0].auto_mount='1'
set fstab.@global[0].delay_root='5'
set fstab.@global[0].check_fs='1'

delete fstab.@mount[0]
add fstab mount
set fstab.@mount[0]=mount
set fstab.@mount[0].enabled='0'
set fstab.@mount[0].fstype='ntfs-3g'

delete fstab.@mount[1]
add fstab mount
set fstab.@mount[1]=mount
set fstab.@mount[1].enabled='0'
set fstab.@mount[1].target='/tmp/plexmediaserver'
set fstab.@mount[1].fstype='squashfs'

delete fstab.@swap[0]
add fstab swap
set fstab.@swap[0]=swap
set fstab.@swap[0].enabled='1'
set fstab.@swap[0].device='/dev/zram0'

delete fstab.@uvol[0]
add fstab uvol
set fstab.@uvol[0]=uvol
set fstab.@uvol[0].initialized='1'
EOF
uci commit fstab
}
exit 0
