# WRX36/bin/extra/files/etc/uci-defaults/z8-luci_statistics
#!/bin/sh

{ [ -f /etc/config/luci_statistics ] && sed s/'[[:space:]]'//g /etc/config/luci_statistics | grep -qE '.+' && [ -f /rom/etc/config/luci_statistics ] && sed s/'[[:space:]]'//g /rom/etc/config/luci_statistics | grep -qE '.+' && ! [ "$(cat /etc/config/luci_statistics)" = "$(cat /rom/etc/config/luci_statistics)" ]; } || {

touch /etc/config/luci_statistics

uci batch <<'EOF'
delete luci_statistics.collectd
set luci_statistics.collectd=statistics
set luci_statistics.collectd.BaseDir='/var/run/collectd'
set luci_statistics.collectd.PIDFile='/var/run/collectd.pid'
set luci_statistics.collectd.PluginDir='/usr/lib/collectd'
set luci_statistics.collectd.TypesDB='/usr/share/collectd/types.db'
set luci_statistics.collectd.Interval='30'
set luci_statistics.collectd.ReadThreads='2'
delete luci_statistics.rrdtool
set luci_statistics.rrdtool=statistics
set luci_statistics.rrdtool.default_timespan='2hour'
set luci_statistics.rrdtool.image_width='600'
set luci_statistics.rrdtool.image_height='150'
set luci_statistics.rrdtool.image_path='/tmp/rrdimg'
delete luci_statistics.collectd_rrdtool
set luci_statistics.collectd_rrdtool=statistics
set luci_statistics.collectd_rrdtool.enable='1'
set luci_statistics.collectd_rrdtool.DataDir='/tmp/rrd'
set luci_statistics.collectd_rrdtool.RRARows='288'
set luci_statistics.collectd_rrdtool.RRASingle='1'
set luci_statistics.collectd_rrdtool.RRATimespans='2hour 1day 1week 1month 1year'
set luci_statistics.collectd_rrdtool.backup='0'
delete luci_statistics.collectd_csv
set luci_statistics.collectd_csv=statistics
set luci_statistics.collectd_csv.enable='0'
set luci_statistics.collectd_csv.StoreRates='0'
set luci_statistics.collectd_csv.DataDir='/tmp'
delete luci_statistics.collectd_email
set luci_statistics.collectd_email=statistics
set luci_statistics.collectd_email.enable='0'
set luci_statistics.collectd_email.SocketFile='/var/run/collectd/email.sock'
set luci_statistics.collectd_email.SocketGroup='nogroup'
delete luci_statistics.collectd_logfile
set luci_statistics.collectd_logfile=statistics
set luci_statistics.collectd_logfile.enable='0'
set luci_statistics.collectd_logfile.LogLevel='notice'
set luci_statistics.collectd_logfile.File='/var/log/collectd.log'
set luci_statistics.collectd_logfile.Timestamp='1'
delete luci_statistics.collectd_network
set luci_statistics.collectd_network=statistics
set luci_statistics.collectd_network.enable='0'
delete luci_statistics.collectd_syslog
set luci_statistics.collectd_syslog=statistics
set luci_statistics.collectd_syslog.enable='0'
set luci_statistics.collectd_syslog.LogLevel='warning'
set luci_statistics.collectd_syslog.NotifyLevel='WARNING'
delete luci_statistics.collectd_unixsock
set luci_statistics.collectd_unixsock=statistics
set luci_statistics.collectd_unixsock.enable='0'
set luci_statistics.collectd_unixsock.SocketFile='/var/run/collectd/query.sock'
set luci_statistics.collectd_unixsock.SocketGroup='nogroup'
delete luci_statistics.collectd_apcups
set luci_statistics.collectd_apcups=statistics
set luci_statistics.collectd_apcups.enable='0'
set luci_statistics.collectd_apcups.Host='localhost'
set luci_statistics.collectd_apcups.Port='3551'
delete luci_statistics.collectd_chrony
set luci_statistics.collectd_chrony=statistics
set luci_statistics.collectd_chrony.enable='0'
set luci_statistics.collectd_chrony.Host='127.0.0.1'
set luci_statistics.collectd_chrony.Port='323'
set luci_statistics.collectd_chrony.Timeout='2'
delete luci_statistics.collectd_conntrack
set luci_statistics.collectd_conntrack=statistics
set luci_statistics.collectd_conntrack.enable='0'
delete luci_statistics.collectd_contextswitch
set luci_statistics.collectd_contextswitch=statistics
set luci_statistics.collectd_contextswitch.enable='0'
delete luci_statistics.collectd_cpu
set luci_statistics.collectd_cpu=statistics
set luci_statistics.collectd_cpu.enable='1'
set luci_statistics.collectd_cpu.ReportByCpu='1'
set luci_statistics.collectd_cpu.ReportByState='1'
set luci_statistics.collectd_cpu.ShowIdle='0'
set luci_statistics.collectd_cpu.ValuesPercentage='1'
delete luci_statistics.collectd_cpufreq
set luci_statistics.collectd_cpufreq=statistics
set luci_statistics.collectd_cpufreq.enable='0'
delete luci_statistics.collectd_curl
set luci_statistics.collectd_curl=statistics
set luci_statistics.collectd_curl.enable='0'
delete luci_statistics.collectd_df
set luci_statistics.collectd_df=statistics
set luci_statistics.collectd_df.enable='0'
set luci_statistics.collectd_df.Devices='/dev/mtdblock/18'
set luci_statistics.collectd_df.MountPoints='/overlay'
set luci_statistics.collectd_df.FSTypes='tmpfs'
set luci_statistics.collectd_df.IgnoreSelected='0'
set luci_statistics.collectd_df.ValuesPercentage='0'
delete luci_statistics.collectd_dhcpleases
set luci_statistics.collectd_dhcpleases=statistics
set luci_statistics.collectd_dhcpleases.enable='0'
set luci_statistics.collectd_dhcpleases.Path='/tmp/dhcp.leases'
delete luci_statistics.collectd_disk
set luci_statistics.collectd_disk=statistics
set luci_statistics.collectd_disk.enable='0'
set luci_statistics.collectd_disk.Disks='hda1 hdb'
set luci_statistics.collectd_disk.IgnoreSelected='0'
delete luci_statistics.collectd_dns
set luci_statistics.collectd_dns=statistics
set luci_statistics.collectd_dns.enable='0'
set luci_statistics.collectd_dns.Interfaces='br-lan'
set luci_statistics.collectd_dns.IgnoreSources='127.0.0.1'
delete luci_statistics.collectd_entropy
set luci_statistics.collectd_entropy=statistics
set luci_statistics.collectd_entropy.enable='0'
delete luci_statistics.collectd_exec
set luci_statistics.collectd_exec=statistics
set luci_statistics.collectd_exec.enable='0'
delete luci_statistics.collectd_interface
set luci_statistics.collectd_interface=statistics
set luci_statistics.collectd_interface.enable='1'
set luci_statistics.collectd_interface.Interfaces='br-lan'
set luci_statistics.collectd_interface.IgnoreSelected='0'
delete luci_statistics.collectd_ipstatistics
set luci_statistics.collectd_ipstatistics=statistics
set luci_statistics.collectd_ipstatistics.enable='0'
delete luci_statistics.collectd_iptables
set luci_statistics.collectd_iptables=statistics
set luci_statistics.collectd_iptables.enable='0'
delete luci_statistics.@collectd_iptables_match[0]
add luci_statistics collectd_iptables_match
set luci_statistics.@collectd_iptables_match[0]=collectd_iptables_match
set luci_statistics.@collectd_iptables_match[0].table='nat'
set luci_statistics.@collectd_iptables_match[0].chain='luci_fw_postrouting'
set luci_statistics.@collectd_iptables_match[0].target='MASQUERADE'
set luci_statistics.@collectd_iptables_match[0].source='10.0.0.0/24'
set luci_statistics.@collectd_iptables_match[0].outputif='br-ff'
set luci_statistics.@collectd_iptables_match[0].name='LAN-Clients traffic'
delete luci_statistics.@collectd_iptables_match[1]
add luci_statistics collectd_iptables_match
set luci_statistics.@collectd_iptables_match[1]=collectd_iptables_match
set luci_statistics.@collectd_iptables_match[1].chain='luci_fw_postrouting'
set luci_statistics.@collectd_iptables_match[1].table='nat'
set luci_statistics.@collectd_iptables_match[1].target='MASQUERADE'
set luci_statistics.@collectd_iptables_match[1].source='10.61.230.0/24'
set luci_statistics.@collectd_iptables_match[1].outputif='br-ff'
set luci_statistics.@collectd_iptables_match[1].name='WLAN-Clients traffic'
delete luci_statistics.collectd_irq
set luci_statistics.collectd_irq=statistics
set luci_statistics.collectd_irq.enable='0'
set luci_statistics.collectd_irq.Irqs='2 3 4 7'
delete luci_statistics.collectd_iwinfo
set luci_statistics.collectd_iwinfo=statistics
set luci_statistics.collectd_iwinfo.enable='1'
delete luci_statistics.collectd_load
set luci_statistics.collectd_load=statistics
set luci_statistics.collectd_load.enable='1'
delete luci_statistics.collectd_memory
set luci_statistics.collectd_memory=statistics
set luci_statistics.collectd_memory.enable='1'
set luci_statistics.collectd_memory.HideFree='0'
set luci_statistics.collectd_memory.ValuesAbsolute='1'
set luci_statistics.collectd_memory.ValuesPercentage='0'
delete luci_statistics.collectd_netlink
set luci_statistics.collectd_netlink=statistics
set luci_statistics.collectd_netlink.enable='0'
set luci_statistics.collectd_netlink.IgnoreSelected='0'
set luci_statistics.collectd_netlink.VerboseInterfaces='br-lan'
set luci_statistics.collectd_netlink.QDiscs='br-lan'
delete luci_statistics.collectd_nut
set luci_statistics.collectd_nut=statistics
set luci_statistics.collectd_nut.enable='0'
set luci_statistics.collectd_nut.UPS='myupsname'
delete luci_statistics.collectd_olsrd
set luci_statistics.collectd_olsrd=statistics
set luci_statistics.collectd_olsrd.enable='0'
set luci_statistics.collectd_olsrd.Port='2006'
set luci_statistics.collectd_olsrd.Host='127.0.0.1'
delete luci_statistics.collectd_openvpn
set luci_statistics.collectd_openvpn=statistics
set luci_statistics.collectd_openvpn.enable='0'
delete luci_statistics.collectd_ping
set luci_statistics.collectd_ping=statistics
set luci_statistics.collectd_ping.enable='0'
set luci_statistics.collectd_ping.TTL='127'
set luci_statistics.collectd_ping.Interval='30'
set luci_statistics.collectd_ping.Hosts='127.0.0.1'
delete luci_statistics.collectd_processes
set luci_statistics.collectd_processes=statistics
set luci_statistics.collectd_processes.enable='0'
set luci_statistics.collectd_processes.Processes='uhttpd dnsmasq dropbear'
delete luci_statistics.collectd_sensors
set luci_statistics.collectd_sensors=statistics
set luci_statistics.collectd_sensors.enable='0'
delete luci_statistics.collectd_snmp6
set luci_statistics.collectd_snmp6=statistics
set luci_statistics.collectd_snmp6.enable='0'
set luci_statistics.collectd_snmp6.Interfaces='br-lan'
set luci_statistics.collectd_snmp6.IgnoreSelected='0'
delete luci_statistics.collectd_splash_leases
set luci_statistics.collectd_splash_leases=statistics
set luci_statistics.collectd_splash_leases.enable='0'
delete luci_statistics.collectd_tcpconns
set luci_statistics.collectd_tcpconns=statistics
set luci_statistics.collectd_tcpconns.enable='0'
set luci_statistics.collectd_tcpconns.ListeningPorts='0'
set luci_statistics.collectd_tcpconns.LocalPorts='22'
set luci_statistics.collectd_tcpconns.AllPortsSummary='0'
delete luci_statistics.collectd_thermal
set luci_statistics.collectd_thermal=statistics
set luci_statistics.collectd_thermal.enable='0'
set luci_statistics.collectd_thermal.IgnoreSelected='0'
delete luci_statistics.collectd_uptime
set luci_statistics.collectd_uptime=statistics
set luci_statistics.collectd_uptime.enable='0'
EOF
uci commit luci_statistics
}
exit 0
