# WRX36/bin/extra/files/etc/uci-defaults/z8-sqm
#!/bin/sh

grep -q act_mirred </etc/modules.d/70-sched-core && sed -i s/'act_mirred'/'act_nssmirred'/ /etc/modules.d/70-sched-core
grep -q qca-nss-match </etc/modules.d/70-sched-core || echo 'qca-nss-match' >> /etc/modules.d/70-sched-core

{ [ -f /etc/config/sqm ] && sed s/'[[:space:]]'//g /etc/config/sqm | grep -qE '.+' && [ -f /rom/etc/config/sqm ] && sed s/'[[:space:]]'//g /rom/etc/config/sqm | grep -qE '.+' && ! [ "$(cat /etc/config/sqm)" = "$(cat /rom/etc/config/sqm)" ]; } || {

touch /etc/config/sqm

uci batch <<'EOF'
delete sqm.eth1
delete sqm.wan
set sqm.wan=queue
set sqm.wan.enabled='0'
set sqm.wan.interface='wan'
set sqm.wan.qdisc='fq_codel'
set sqm.wan.script='nss-zk.qos'
set sqm.wan.overhead='44'
set sqm.wan.linklayer='ethernet'
set sqm.wan.linklayer_advanced='1'
set sqm.wan.tcMPU='84'
set sqm.wan.ingress_ecn='ECN'
set sqm.wan.egress_ecn='NOECN'
set sqm.wan.squash_dscp='0'
set sqm.wan.squash_ingress='0'
set sqm.wan.tcMTU='2047'
set sqm.wan.tcTSIZE='512'
set sqm.wan.qdisc_advanced='1'
set sqm.wan.qdisc_really_really_advanced='1'
set sqm.wan.itarget='auto'
set sqm.wan.etarget='auto'
EOF
uci commit sqm
}
exit 0
