# WRX36/bin/extra/files/etc/uci-defaults/z8-umdns
#!/bin/sh

{ [ -f /etc/config/umdns ] && sed s/'[[:space:]]'//g /etc/config/umdns | grep -qE '.+' && [ -f /rom/etc/config/umdns ] && sed s/'[[:space:]]'//g /rom/etc/config/umdns | grep -qE '.+' && ! [ "$(cat /etc/config/umdns)" = "$(cat /rom/etc/config/umdns)" ]; } || {

touch /etc/config/umdns

uci batch <<'EOF'
delete umdns.@umdns[0]
add umdns umdns
set umdns.@umdns[0]=umdns
set umdns.@umdns[0].jail='1'
add_list umdns.@umdns[0].network='lan'
add_list umdns.@umdns[0].network='IoT'
EOF
uci commit umdns
}
exit 0
