# WRX36/bin/extra/files/etc/uci-defaults/z9-done
#!/bin/sh
/etc/init.d/network restart

{ sleep 10; reboot; } &

exit 0
