# WRX36/bin/extra/files/etc/uci-defaults/z9-done
#!/bin/sh
/etc/init.d/network restart
exit 0
